# app-backup/backrest/files/backrest.initd
#!/sbin/openrc-run

name="backrest"
description="backrest Daemon"
command="/usr/bin/backrest"
command_args="-data-dir /var/lib/backrest/ -config-file /var/lib/backrest/config.json"
command_background=true
pidfile="/run/backrest.pid"
supervisor=supervise-daemon

output_log="/var/log/backrest/backrest.log"
error_log="/var/log/backrest/backrest.log"

depend() {
    need net
}
