# app-misc/zigbee2mqtt/files/zigbee2mqtt
#!/sbin/openrc-run

depend() {
	need net
}

checkconfig() {
	if [ ! -f /var/lib/${SVCNAME}/configuration.yaml ] ; then
		eerror "No /var/lib/${SVCNAME}/configuration.yaml file exists!"
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --env ZIGBEE2MQTT_DATA=/var/lib/${SVCNAME} --env  --chdir /opt/${SVCNAME} --exec npm start --background --make-pidfile --pidfile /var/run/${SVCNAME}.pid
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --pidfile /var/run/${SVCNAME}.pid
	eend $?
}
